# bspwm/.config/systemd/user/battery-notify.service
[Unit]
Description=Notify when battery is low

[Service]
Type=oneshot
ExecStart=/home/jay/dotfiles/scripts/bats.sh

